# packaging/systemd/tailnetd.service
# tailnetd.service — systemd unit for the tailscaled-rs daemon (Linux).
#
# Runs `tailnetd` as an always-on SYSTEM service (root-owned, started at boot via
# multi-user.target). This is the system-wide deployment; a per-user-session variant is
# possible later by shipping a `--user` unit (placed under ~/.config/systemd/user/ and
# managed with `systemctl --user`), which would drop the system hardening below and resolve
# the state dir from $XDG_STATE_HOME instead of /var/lib.
#
# This unit is HARDENED on purpose: tailnetd holds unencrypted WireGuard/node key material in
# its state dir, so the service is locked down with the systemd sandboxing directives below.
#
# WARNING: experimental, unaudited software — not for production. See the repository README.
# Install / uninstall instructions: ../README.md

[Unit]
Description=tailscaled-rs daemon (tailnetd) — experimental WireGuard mesh node
Documentation=https://github.com/GeiserX/tailscaled-rs
# Wait for routable network before starting (the node dials a control server + DERP).
After=network-online.target
Wants=network-online.target

[Service]
# Type=notify (matches Go's tailscaled): tailnetd sends sd_notify(READY=1) once the LocalAPI socket is
# bound and the backend has loaded + attempted auto-start (i.e. it can actually serve), so systemd
# holds `systemctl start` until the daemon is genuinely ready rather than merely exec'd. The daemon
# reads $NOTIFY_SOCKET and emits the readiness datagram itself (no libsystemd dependency); run outside
# systemd (NOTIFY_SOCKET unset) the notify is a harmless no-op.
Type=notify

# Operator tunables, optional. Go ships EnvironmentFile=/etc/default/tailscaled; we use the same seam
# so an admin can set env knobs (e.g. TAILNETD_LOG, or extra engine env) WITHOUT editing this vendor
# unit (which an upgrade overwrites). The leading `-` makes it optional: a system without the file
# starts fine. It is read BEFORE the hardcoded Environment= lines below, so those win — the file
# CANNOT relocate the state dir (StateDirectory= + ReadWritePaths= are fixed to /var/lib/tailnetd) or
# disable the experiment gate; it only adds/overrides knobs that do not fight the sandbox.
EnvironmentFile=-/etc/default/tailnetd

ExecStart=/usr/local/bin/tailnetd
# Revert OS-level state when the service stops OR crashes (Go: ExecStopPost=tailscaled --cleanup).
# `--cleanup` removes the daemon's own stale LocalAPI socket and is structurally incapable of touching
# the node key or prefs.json (it is handed only the socket path), so a crash/stop never destroys
# identity — it just clears the transient socket so the next start binds cleanly. Runs with this
# unit's environment, so it inherits TAILNETD_STATE_DIR and resolves the same <statedir>/tailnetd.sock.
ExecStopPost=/usr/local/bin/tailnetd --cleanup
Restart=on-failure
RestartSec=5

# --- Experiment opt-in (REQUIRED) ---------------------------------------------------------
# The underlying engine refuses to start unless this is set to exactly this value. We set it
# here, explicitly and visibly, because enabling this unit IS the operator opting in to running
# unaudited, experimental software on purpose. tailnetd deliberately does NOT set it for itself.
Environment=TS_RS_EXPERIMENT=this_is_unstable_software

# --- State directory ----------------------------------------------------------------------
# systemd creates /var/lib/tailnetd with 0700 ownership root:root (StateDirectory). The daemon
# additionally enforces 0700 on this dir itself at startup (it holds unencrypted key material),
# so the two layers agree. TAILNETD_STATE_DIR points the daemon at this exact path.
StateDirectory=tailnetd
StateDirectoryMode=0700
Environment=TAILNETD_STATE_DIR=/var/lib/tailnetd

# --- Hardening (this daemon holds private keys) -------------------------------------------
# Drop ALL privileges and sandbox aggressively. The MVP is userspace-networking only (DERP /
# in-process WireGuard), so it needs no capabilities, no TUN device, and no write access to the
# system outside its own state dir.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true
SystemCallArchitectures=native
# ProtectSystem=strict makes the whole filesystem read-only; carve out only the state dir.
ReadWritePaths=/var/lib/tailnetd
# Empty set = strip every capability. The userspace-networking MVP needs none.
CapabilityBoundingSet=

# --- Syscall / namespace lockdown ---------------------------------------------------------
# Tighten the file-creation mask and the syscall + address-family surface for this key-holding
# daemon. UMask=0077 is the HIGHEST (most restrictive) value: the engine writes node.key.json
# with the inherited umask, so this closes the brief world-/group-readable window on the key
# file even though it lives inside the 0700 state dir.
UMask=0077
# AF_UNIX = LocalAPI socket; AF_INET/AF_INET6 = UDP + DERP; AF_NETLINK = the engine's interface
# queries (do NOT drop netlink — removing it can break the engine's address enumeration).
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
# Allow only the generic service syscall set, then subtract privileged + resource-control groups.
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM
# Hide other processes and expose only pid entries under /proc; deny device-node access entirely
# (the userspace MVP opens no /dev/net/tun).
ProtectProc=invisible
ProcSubset=pid
PrivateDevices=true

# --- Phase-3 (TUN / kernel-mode) NOTE -----------------------------------------------------
# When the kernel TUN data path lands (Phase 3), this unit must be RELAXED to grant the daemon
# the ability to create and configure a tun interface. At minimum you will then need:
#   CapabilityBoundingSet=CAP_NET_ADMIN
#   AmbientCapabilities=CAP_NET_ADMIN
#   DeviceAllow=/dev/net/tun rw
#   PrivateDevices=false   # PrivateDevices=true hides /dev/net/tun; TUN mode needs it visible
#   # and loosen ProtectSystem (e.g. ProtectSystem=full) so /dev/net/tun is reachable,
#   # plus relax RestrictNamespaces if a network namespace is used. You will also likely need
#   # to widen the syscall/address-family lockdown above for the tun ioctls: relax
#   # SystemCallFilter (the @privileged subtraction blocks some interface-config calls) and
#   # extend RestrictAddressFamilies (e.g. AF_PACKET) if the kernel data path requires it.
# Do NOT add these now: the current MVP does not open /dev/net/tun, and granting CAP_NET_ADMIN
# to experimental code with no need for it only widens the blast radius.

[Install]
WantedBy=multi-user.target
